# config/launchd/jp.radio-calisthenics-together.prepare.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>jp.radio-calisthenics-together.prepare</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>scripts/prepare_environment.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>50</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{REPO_DIR}}/logs/prepare_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{REPO_DIR}}/logs/prepare_stderr.log</string>
    <key>WorkingDirectory</key>
    <string>{{REPO_DIR}}</string>
</dict>
</plist>
